# compilators/debian/nginx/postrm
#!/bin/sh

update-rc.d phyre remove > /dev/null 2>&1

exit 0
